# gsb4-deb/coccoc-gsb4.cron.d
MAILTO=quannk@coccoc.com
LOGDIR=/var/log/coccoc-gsb4

*/30 * * * * root gsb4.sh
